# ttl/start_ttl_server.sh
# /bin/sh
# chkconfig: 2345 80 90
# description: satrt ttl-server 

/home/pi/MiniAnaconda3/bin/python /home/pi/Documents/code/ttl/ttl-server-update.py
